# server/alpine/cloudflared
#!/sbin/openrc-run

name="cloudflared"
description="Cloudflare Tunnel"

command="/usr/bin/cloudflared"
command_args="tunnel --config /root/.cloudflared/config.yml run"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"
output_log="/var/log/cloudflared.log"
error_log="/var/log/cloudflared.log"

depend() {
    need net
    after coloringbook
}

start_pre() {
    checkpath --file --mode 0644 "$output_log"
}
